# dist/etc/systemd/system/localrepo.service
[Unit]
Description=Sync local pacman repo from S3

[Service]
Type=oneshot
WorkingDirectory=@@REPO_DIR@@
ExecStart=/usr/bin/just pull
